# scrapyd_logparser_sup/supervisord.conf
[supervisord]
nodaemon=true

[program:scrapyd]
command= scrapyd
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:logparser]
command= logparser
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
